# luci-app-localmachines/root/etc/uci-defaults/40_localmachines
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.localmachines
	set ucitrack.localmachines=localmachines
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
